# dz-06/01file/searchlog.timer
[Unit]
Description=DZ_06

[Timer]
OnUnitActiveSec=30s

Unit=searchlog.service

[Install]
WantedBy=timers.target
