# dz-06/01file/searchlog.service
[Unit]
Description=DZ_06

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/searchlog
ExecStart=/bin/grep $STR $FILE
